# provision/etc/profile.d/furiend.sh
alias ll="ls -Alh --group-directories-first --color=always"

export PS1="$(tput setaf 5)\h$(tput sgr0) \w $(tput setaf 5)#$(tput sgr0) "
export PATH="$PATH:/furiend/bin"
